# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Comment[tr]=Renk sıcaklığı ayarlama aracı
GenericName[zh_CN]=色温调节
Name[fr]=Redshift
Name=Redshift
Comment[da]=Justeringsværktøj for farvetemperatur
Comment[it]=Strumento per la regolazione della temperatura del colore
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
GenericName[it]=Regolazione della temperatura del colore
Name[de]=Redshift
Comment[ru]=Инструмент регулирования цветовой температуры
Name[zh_CN]=红移
GenericName[tr]=Renk sıcaklığı ayarı
Hidden=true
Categories=Utility;
GenericName=Color temperature adjustment
Name[hu]=Redshift
GenericName[es]=Ajuste de la temperatura de color
GenericName[lt]=Spalvos temperatūros reguliavimas
Comment[zh_CN]=色温调节工具
Comment[fr]=Outil de réglage de la température de couleur
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Comment[es]=Herramienta para el ajuste de la temperatura de color
Icon=redshift
X-GNOME-Autostart-enabled=false
Comment=Color temperature adjustment tool
GenericName[de]=Farbtemperatur Anpassung
Exec=redshift-gtk
Type=Application
Name[es]=Redshift
Name[it]=Redshift
GenericName[da]=Justering af farvetemperatur
Name[pt_BR]=Redshift
Comment[de]=Farbtemperatur Anpassungswerkzeug
Terminal=false
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
GenericName[pl]=Dostosowanie temperatury barwowej
Comment[hu]=Színhőmérséklet beállító eszköz
GenericName[hu]=Szín hőmérséklet beállítás
Name[ru]=Redshift
StartupNotify=true
Name[pl]=Redshift
Name[lt]=Redshift
GenericName[cs]=Přizpůsobení barevné teploty
GenericName[fr]=Réglage de la température de couleur
Version=1.0
GenericName[nl]=Kleurtemperatuur aanpassing
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
Name[tr]=Redshift
GenericName[ru]=Настройка цветовой температуры

